# packaging/com.dockswitch.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.dockswitch</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/DockSwitch</string>
        <string>__HOME__/.config/dockswitch/config.json</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/dockswitch/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/dockswitch/stderr.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>

</dict>
</plist>
